# buildroot/board/rockchip/common/base/etc/udev/rules.d/99-rpdzkj-mount-M2.rules
#!/bin/sh
ACTION=="add", KERNEL=="nvme0n[0-9]p[0-9]",RUN+="/etc/mount-M2.sh %k"
